# packages_dpkg/wireplumber_pkg/postinst
#!/bin/bash

systemctl --global enable wireplumber.service

