# packages_dpkg/pipewire/postinst
#!/bin/bash


systemctl --user enable pipewire

systemctl --user enable  pipewire-pulse

systemctl --user enable wireplumber
