# packages/basefile/initd/sysklogd
#!/bin/sh

### BEGIN INIT INFO
# Provides:            $syslog
# Required-Start:      $first localnet
# Should-Start:
# Required-Stop:       $local_fs
# Should-Stop:         sendsignals
# Default-Start:       2 3 4 5
# Default-Stop:        0 1 6
# Short-Description:   Starts kernel and system log daemons.
# Description:         Starts kernel and system log daemons.
#                      /etc/fstab.
### END INIT INFO

. /etc/init.d/function

case "${1}" in
   start)
      log_info_msg "Starting system log daemon..."
      parms=${SYSKLOGD_PARMS-'-m 0'}
      start_daemon /usr/sbin/syslogd $parms
      evaluate_retval

      log_info_msg "Starting kernel log daemon..."
      start_daemon /usr/sbin/klogd
      evaluate_retval
      ;;

   stop)
      log_info_msg "Stopping kernel log daemon..."
      killproc /usr/sbin/klogd
      evaluate_retval

      log_info_msg "Stopping system log daemon..."
      killproc /usr/sbin/syslogd
      evaluate_retval
      ;;

   reload)
      log_info_msg "Reloading system log daemon config file..."
      pid=`pidofproc syslogd`
      kill -HUP "${pid}"
      evaluate_retval
      ;;

   restart)
      ${0} stop
      sleep 1
      ${0} start
      ;;

   status)
      statusproc /usr/sbin/syslogd
      statusproc klogd
      ;;
   enable)
      mkdir -p /etc/rc.d/rcS.d
      cd /etc/rc.d/rcS.d && ln -sf ../../init.d/sysklogd S03sysklogd
      log_info_msg "Service sysklogd enable...\n\n"
      ;;

   disable)
      rm /etc/rc.d/rcS.d/S03sysklogd
      log_info_msg "Service sysklogd enable...\n\n"
      ;;
   *)
      echo "Usage: ${0} {start|stop|reload|restart|status|enable|disable}"
      exit 1
      ;;
esac

exit 0

# End sysklogd
