# packages_dpkg/make-ca/postinst
#!/bin/bash


/usr/sbin/make-ca -g
systemctl enable update-pki.timer
