# packages/irbis-bootscripts/bootscript/initd/dbus
#!/bin/sh

### BEGIN INIT INFO
# Provides:            dbus messagebus
# Required-Start:
# Should-Start:        $syslog random
# Required-Stop:
# Should-Stop:         $syslog
# Default-Start:       2 3 4 5
# Default-Stop:        0 1 6
# Short-Description:   Starts message bus.
# Description:         Starts message bus.
### END INIT INFO

. /etc/init.d/function

#$LastChangedBy$
#$Date$

pidfile=/run/dbus/pid
socket=/run/dbus/system_bus_socket

case "$1" in
   start)
      log_info_msg "Starting the D-Bus Messagebus Daemon..."
      mkdir -p /run/dbus
      /usr/bin/dbus-uuidgen --ensure
      start_daemon /usr/bin/dbus-daemon --system
      evaluate_retval
      ;;

   stop)
      log_info_msg "Stopping the D-Bus Messagebus Daemon..."
      killproc /usr/bin/dbus-daemon
      evaluate_retval
      rm -f $socket $pidfile
      ;;

   restart)
      $0 stop
      sleep 1
      $0 start
      ;;

   status)
      statusproc /usr/bin/dbus-daemon
      ;;
   enable)
        mkdir -p ../rc.d/rcS.d/
   	cd ../rc.d/rcS.d/ && ln -sf ../../init.d/dbus S01dbus
   	log_info_msg "Service enable..."
      ;;
   disable)
   	rm ../rc.d/rcS.d/S01dbus
   	log_info_msg "Service disable..."
      ;;   
   *)
      echo "Usage: $0 {start|stop|restart|status|enable|disable}"
      exit 1
      ;;
esac

# End /etc/init.d/dbus
